# deploy/systemd/financeiro-restart.service
[Unit]
Description=Reinício agendado do backend Django (financeiro)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart financeiro.service
